# supervisord.conf
[supervisord]
nodaemon=true

[program:netuitive-statsd]
command=python /opt/netuitive/netuitive-statsd/netuitive-statsd.py -d -c /etc/netuitive/netuitive-statsd/netuitive-agent.conf start 

